# dockerconf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/sbin/nginx
autostart=true
autorestart=true
